# apps/cli/templates/systemd-user/codex-swarm-worker@.service
[Unit]
Description=Codex Swarm Local Worker Daemon %i (single-host deployment)
After=network-online.target codex-swarm-api.service
Requires=codex-swarm-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_ROOT__
EnvironmentFile=__ENV_FILE__
EnvironmentFile=-%h/.config/codex-swarm/workers/%i.env
Restart=always
RestartSec=5
ExecStartPre=/bin/sh -lc 'mkdir -p "$CODEX_SWARM_WORKSPACE_ROOT"'
ExecStart=/bin/sh -lc 'export CODEX_SWARM_SERVICE_TOKEN="${CODEX_SWARM_SERVICE_TOKEN:-${AUTH_SERVICE_TOKEN:?missing CODEX_SWARM_SERVICE_TOKEN or AUTH_SERVICE_TOKEN}}" CODEX_SWARM_SERVICE_NAME="${CODEX_SWARM_SERVICE_NAME:-local-daemon}"; node apps/api/dist/src/ops/local-worker-daemon.js'

[Install]
WantedBy=default.target
