# apps/cli/templates/systemd-user/codex-swarm-api.service
[Unit]
Description=Codex Swarm API (single-host deployment)
After=network-online.target codex-swarm-postgres.service codex-swarm-redis.service
Requires=codex-swarm-postgres.service codex-swarm-redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_ROOT__
EnvironmentFile=__ENV_FILE__
Restart=always
RestartSec=5
ExecStartPre=/bin/sh -lc 'export CODEX_SWARM_API_BASE_URL="${CODEX_SWARM_API_BASE_URL:-http://${CODEX_SWARM_HOST}:${CODEX_SWARM_API_PORT}}" VITE_API_BASE_URL="${CODEX_SWARM_API_BASE_URL:-http://${CODEX_SWARM_HOST}:${CODEX_SWARM_API_PORT}}"; node ops/deploy/write-frontend-runtime-config.mjs'
ExecStartPre=/bin/sh -lc 'export NODE_ENV=production HOST="$CODEX_SWARM_HOST" PORT="$CODEX_SWARM_API_PORT" DATABASE_URL="$CODEX_SWARM_DATABASE_URL" ARTIFACT_STORAGE_ROOT="$CODEX_SWARM_ARTIFACT_STORAGE_ROOT" ARTIFACT_BASE_URL="$CODEX_SWARM_ARTIFACT_BASE_URL" DEV_AUTH_TOKEN="$CODEX_SWARM_DEV_AUTH_TOKEN"; node apps/api/dist/src/db/migrate.js'
ExecStart=/bin/sh -lc 'export NODE_ENV=production HOST="$CODEX_SWARM_HOST" PORT="$CODEX_SWARM_API_PORT" FRONTEND_DIST_ROOT="frontend/dist" DATABASE_URL="$CODEX_SWARM_DATABASE_URL" ARTIFACT_STORAGE_ROOT="$CODEX_SWARM_ARTIFACT_STORAGE_ROOT" ARTIFACT_BASE_URL="$CODEX_SWARM_ARTIFACT_BASE_URL" DEV_AUTH_TOKEN="$CODEX_SWARM_DEV_AUTH_TOKEN"; node apps/api/dist/src/server.js'

[Install]
WantedBy=default.target
