# ops/systemd-user/codex-swarm.target
[Unit]
Description=Codex Swarm Tailnet-Only Stack
Wants=codex-swarm-postgres.service codex-swarm-redis.service codex-swarm-api.service codex-swarm-frontend.service codex-swarm-worker.service codex-swarm-worker@worker-2.service codex-swarm-worker@worker-3.service codex-swarm-worker@worker-4.service
