# ops/systemd-user/codex-swarm-worker.service
[Unit]
Description=Codex Swarm Local Worker Daemon (tailnet-only deployment)
After=network-online.target codex-swarm-api.service
Requires=codex-swarm-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/florian/codex-swarm
EnvironmentFile=%h/.config/codex-swarm/tailnet.env
Restart=always
RestartSec=5
ExecStartPre=/bin/sh -lc 'corepack pnpm install --frozen-lockfile'
ExecStartPre=/bin/sh -lc 'mkdir -p "$CODEX_SWARM_WORKSPACE_ROOT"'
ExecStart=/bin/sh -lc 'export CODEX_SWARM_SERVICE_TOKEN="${CODEX_SWARM_SERVICE_TOKEN:-${AUTH_SERVICE_TOKEN:?missing CODEX_SWARM_SERVICE_TOKEN or AUTH_SERVICE_TOKEN}}" CODEX_SWARM_SERVICE_NAME="${CODEX_SWARM_SERVICE_NAME:-local-daemon}"; corepack pnpm --dir apps/api exec tsx src/ops/local-worker-daemon.ts'

[Install]
WantedBy=default.target
