# ops/systemd-user/codex-swarm-frontend.service
[Unit]
Description=Codex Swarm Frontend (tailnet-only deployment)
After=network-online.target codex-swarm-api.service
Requires=codex-swarm-api.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/florian/codex-swarm
EnvironmentFile=%h/.config/codex-swarm/tailnet.env
Restart=always
RestartSec=5
ExecStartPre=/bin/sh -lc 'corepack pnpm install --frozen-lockfile'
ExecStartPre=/bin/sh -lc 'export VITE_API_BASE_URL="$CODEX_SWARM_API_BASE_URL" VITE_API_TOKEN="$CODEX_SWARM_DEV_AUTH_TOKEN"; corepack pnpm --dir frontend build'
ExecStart=/bin/sh -lc 'corepack pnpm --dir frontend preview --host "$CODEX_SWARM_HOST" --port "$CODEX_SWARM_FRONTEND_PORT"'

[Install]
WantedBy=default.target
